# systemd/.config/systemd/user/cliphist.service
[Unit]
Description=Clipboard history

[Service]
ExecStart=cliphist store

[Install]
WantedBy=graphical-session.target
