# systemd/.config/systemd/user/mood-engine.service
[Unit]
Description=Generate palette
After=wallpaper.target

[Service]
Type=oneshot
WorkingDirectory=%h/hypraura
ExecStart=%h/hypraura/scripts/mood-engine
