# systemd/.config/systemd/user/scenes.timer
[Unit]
Description=Rotate scenes hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=scenes.service

[Install]
WantedBy=timers.target
