# systemd/.config/systemd/user/swaync.service
[Unit]
Description=Notifications

[Service]
ExecStart=swaync

[Install]
WantedBy=graphical-session.target
