# systemd/.config/systemd/user/swww.service
[Unit]
Description=Wallpaper daemon
After=graphical-session.target

[Service]
ExecStart=swww-daemon

[Install]
WantedBy=default.target
