# systemd/.config/systemd/user/theming.target
[Unit]
Description=Apply theming
Requires=mood-engine.service
After=mood-engine.service
